# systemd/mqvpn-server.service
[Unit]
Description=mqvpn VPN Server
After=network-online.target
Wants=network-online.target
Documentation=https://github.com/mp0rta/mqvpn

[Service]
Type=simple
ExecStartPre=/usr/local/lib/mqvpn/mqvpn-server-nat.sh setup /etc/mqvpn/server.conf
ExecStart=/usr/local/bin/mqvpn --config /etc/mqvpn/server.conf
ExecStopPost=/usr/local/lib/mqvpn/mqvpn-server-nat.sh teardown /etc/mqvpn/server.conf
Restart=on-failure
RestartSec=5
StartLimitIntervalSec=60
StartLimitBurst=3

# Security hardening
ProtectHome=true
ProtectSystem=strict
PrivateTmp=true
ReadOnlyPaths=/etc/mqvpn
ReadWritePaths=/dev/net/tun

[Install]
WantedBy=multi-user.target
